# tiny/profiles/bathroom/rootfs/etc/inittab

::sysinit:/etc/init.d/rcS
::respawn:/bin/sh
::restart:/sbin/init
::shutdown:/bin/umount -a -r
